# ME58800/uno1.conf
[program:uno1]
directory=home/tam/me588project
command=bash start-uno1.sh
autostart=true
autorestart=true
user=root