# ABE53000/project/env/camera_interface/camera_backend.conf
[program:camera_backend]
directory=/home/pi/abe591-project/backend
command=bash supervisor.sh
autostart=true
autorestart=true
user=pi

