# ABE53000/project/env/live_stream/live_stream.conf
[program:live_stream]
directory=/home/pi/abe591-project/backend/live_stream
command=bash supervisor.sh
autostart=true
autorestart=true
user=pi

